# components/php54/files/php-fpm.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
	<service name='network/http/php-fpm' type='service' version='0'>
		<dependency name='network' grouping='require_all' restart_on='error' type='service'>
			<service_fmri value='svc:/milestone/network:default'/>
		</dependency>
		<dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
			<service_fmri value='svc:/system/filesystem/local:default'/>
		</dependency>
		<dependent name='apache22' restart_on='none' grouping='optional_all'>
			<service_fmri value='svc:/network/http/apache22:default'/>
		</dependent>
		<dependent name='nginx' restart_on='none' grouping='optional_all'>
			<service_fmri value='svc:/network/http/nginx:default'/>
		</dependent>
		<exec_method name='start' type='method' exec='/ec/lib/php/5.4/bin/php-fpm -c %{php-fpm/php_config} -y %{php-fpm/config_file}' timeout_seconds='60'/>
		<exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
		<exec_method name='refresh' type='method' exec=':kill -HUP' timeout_seconds='60'/>
		<property_group name='php-fpm' type='application'>
			<propval name='config_file' type='astring' value=''/>
			<propval name='php_config' type='astring' value='/ec/etc/php/5.4/php.ini'/>
		</property_group>
		<property_group name='startd' type='framework'>
			<propval name='ignore_error' type='astring' value='core,signal'/>
		</property_group>
		<stability value='Evolving'/>
	</service>
</service_bundle>
