# components/spread/files/spread.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='network/spread' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='network-service' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/service'/>
    </dependency>
    <dependency name='name-services' grouping='require_all' restart_on='refresh' type='service'>
      <service_fmri value='svc:/milestone/name-services'/>
    </dependency>
    <dependency name='network-physical' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/physical'/>
    </dependency>
    <dependency name='network-routing-setup' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/routing-setup:default'/>
    </dependency>
      <method_context>
      <method_credential group='root' user='root'/>
    </method_context>
    <exec_method name='start' type='method' exec='/ec/bin/spread -c /ec/etc/spread.conf' timeout_seconds='10'/>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='10'/>
    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='child'/>
      <propval name='ignore_error' type='astring' value='core,signal'/>
    </property_group>
    <stability value='Stable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Spread Daemon</loctext>
      </common_name>
      <documentation>
        <manpage title='spread' section='1' manpath='/ec/share/man'/>
      </documentation>
    </template>
  </service>
</service_bundle>
