# components/logstash-forwarder/logstash-forwarder.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='system/logstash-forwarder' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name='milestone' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/sysconfig'/>
    </dependency>
    <exec_method name='start' type='method' exec='/ec/bin/logstash-forwarder -config=/ec/etc/logstash-forwarder/logstash-forwarder.conf -quiet=true &amp;' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
    <stability value='Unstable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>logstash-forwarder</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
