# components/core-os/files/var/svc/manifest/network/rpc/100235_1-rpc_ticotsord.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    Service manifest for the 100235 service.

    Generated by inetconv(1M) from inetd.conf(4).
-->

<service_bundle type='manifest' name='inetconv:100235/1'>

<service
	name='network/rpc-100235_1/rpc_ticotsord'
	type='service'
	version='1'>

	<create_default_instance enabled='true'/>

	<restarter>
		<service_fmri value='svc:/network/inetd:default' />
	</restarter>

	<dependency
		name='rpcbind'
		grouping='require_all'
		restart_on='restart'
		type='service'>
		<service_fmri value='svc:/network/rpc/bind' />
	</dependency>

	<!--
	    Set a timeout of 0 to signify to inetd that we don't want to
	    timeout this service, since the forked process is the one that
	    does the service's work. This is the case for most/all legacy
	    inetd services; for services written to take advantage of SMF
	    capabilities, the start method should fork off a process to
	    handle the request and return a success code.
	-->
	<exec_method
		type='method'
		name='inetd_start'
		exec='/usr/lib/fs/cachefs/cachefsd'
		timeout_seconds='0'>
		<method_context>
			<method_credential user='root' group='root' />
		</method_context>
	</exec_method>

	<!--
	    Use inetd's built-in kill support to disable services.
	-->
	<exec_method
		type='method'
		name='inetd_disable'
		exec=':kill'
		timeout_seconds='0'>
	</exec_method>

	<!--
	    Use inetd's built-in process kill support to offline wait type
	    services.
	-->
	<exec_method
		type='method'
		name='inetd_offline'
		exec=':kill_process'
		timeout_seconds='0'>
	</exec_method>

	<!--
	    This property group is used to record information about
	    how this manifest was created.  It is an implementation
	    detail which should not be modified or deleted.
	-->
	<property_group name='inetconv' type='framework'>
		<propval name='converted' type='boolean' value='true' />
		<propval name='version' type='integer' value='1' />
		<propval name='source_line' type='astring' value=
'100235/1 tli rpc/ticotsord wait root /usr/lib/fs/cachefs/cachefsd cachefsd'
		/>
	</property_group>

	<property_group name='inetd' type='framework'>
		<propval name='name' type='astring' value='100235' />
		<propval name='endpoint_type' type='astring' value='tli' />
		<propval name='proto' type='astring' value='ticotsord' />
		<propval name='wait' type='boolean' value='true' />
		<propval name='isrpc' type='boolean' value='true' />
		<propval name='rpc_low_version' type='integer' value='1' />
		<propval name='rpc_high_version' type='integer' value='1' />
	</property_group>

	<stability value='External' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
100235
			</loctext>
		</common_name>
	</template>
</service>

</service_bundle>
