# components/exim/exim-manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
  <service_bundle type='manifest' name='exim'>
  <service name='network/smtp/exim' type='service' version='0'>
    <dependency name='fs' grouping='require_all'  restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local' />
    </dependency>

    <dependency name='net' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default' />
    </dependency>

    <dependency name='ns' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/name-services' />
    </dependency>

    <exec_method type='method' name='start'
      exec='/ec/var/svc/method/exim.sh start'
      timeout_seconds='-1'>
    </exec_method>

    <exec_method type='method' name='stop'
      exec=':kill'
      timeout_seconds='30'>
    </exec_method>

    <exec_method type='method' name='refresh'
      exec=':kill -HUP'
      timeout_seconds='30'>
    </exec_method>

    <instance name='default' enabled='false'>
      <property_group name='exim' type='application'>
        <propval name='startup_options' type='astring' value='-bd -q15m'/>
        <propval name='configuration_file' type='astring' value='/ec/etc/exim/exim.conf'/>
        <propval name='exim_binary' type='astring' value='/ec/sbin/exim'/>
      </property_group>
    </instance>
  </service>
  </service_bundle>

