# components/php72/files/php72-postinstall.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
	Copyright (C) EveryCity 2013
-->
<service_bundle type='manifest' name='php72-postinstall'>

<service
	name='tmp/postinstall/php72'
	type='service'
	version='1'>

	<exec_method
		type='method'
		name='start'
		exec='sh -c &quot;[ -h /ec/lib/php/active ] || ln -s 7.2 /ec/lib/php/active&quot; ; /usr/bin/ctrun svcadm disable tmp/postinstall/php72'
		timeout_seconds='60'/>

 	<exec_method
		type='method'
		name='stop'
		exec='/usr/bin/ctrun svccfg delete -f tmp/postinstall/php72'
		timeout_seconds='30'/>

	<instance name='default' enabled='true'/>

	<stability value='Evolving'/>

</service>

</service_bundle>
